# conf/upstart/eclimd.conf
start on startup
console log
task
exec /usr/local/bin/eclimd &
