# conf/upstart/xscreensaver.conf
start on startup
console log
task
exec /usr/bin/xscreensaver &
